# crontab
@reboot /usr/local/bin/python /app/import_adstat.py >> /var/log/cron.log 2>&1
0 * * * * /usr/local/bin/python /app/import_adstat.py >> /var/log/cron.log 2>&1
